# external/arl_dhs/share_dir.service
[Unit]
Description=Share share dir

[Service]
Type=oneshot
ExecStart=/etc/share_dir.sh

[Install]
WantedBy=multi-user.target
